# apt-cacher/ubuntu/supervisord.conf
[supervisord]
nodaemon=true

[program:apt-cacher]
command = /usr/bin/perl /usr/sbin/apt-cacher -R 3 -d -p /var/run/apt-cacher.pid
